# cron.d/docker
35 02 * * * root /usr/bin/docker system prune -f >/dev/null
40 02 * * * root /opt/adm/bin/registry-cli -d -n 6 --keep-tags-like oldprod prod staging latest >/dev/null
45 02 * * * root /opt/adm/bin/registry-purge >/dev/null 2>&1
